# deploy/systemd/sc-backup.service
[Unit]
Description=SentinelCore encrypted database backup
Documentation=https://github.com/hsynokyay/sentinelcore/blob/main/docs/superpowers/plans/2026-04-18-phase7-data-security.md
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/sentinelcore/env/backup.env

[Service]
Type=oneshot
User=sentinelcore
Group=sentinelcore
# backup.env supplies:
#   DATABASE_URL=postgres://sentinelcore_readonly:...@localhost:5432/sentinelcore
#   SC_BACKUP_AGE_RECIPIENTS=age1recipient1... age1recipient2...
EnvironmentFile=/opt/sentinelcore/env/backup.env

# Readonly role only — backups MUST NOT hold the write role.
# The role's password lives alongside the service creds; backup.env
# overrides DATABASE_URL for this unit only.
ExecStart=/usr/local/bin/sc-backup \
    --out /var/lib/sentinelcore/backups/sc-%%Y-%%m-%%dT%%H-%%M-%%SZ.sql.age

# Sandbox — this unit can read the env file and write backups; nothing else.
ProtectSystem=strict
ReadWritePaths=/var/lib/sentinelcore/backups
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

# The binary runs pg_dump — allow the standard network and filesystem
# syscalls, deny everything else.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Two-hour cap matches the CLI's internal timeout.
TimeoutStartSec=2h

[Install]
WantedBy=multi-user.target
